# config/systemd/user/msu.service
[Unit]
Description=Monitor setup utility

[Service]
Type=oneshot
ExecStart=%h/.local/bin/msu
