# config/systemd/user/polybar.service
[Unit]
Description=Status bar

[Service]
ExecStart=/usr/bin/polybar base
