# config/systemd/user/xautolock.service
[Unit]
Description=Automatic screen locker

[Service]
ExecStart=/usr/bin/xautolock -time 15 -locker slock -secure
