# config/systemd/user/xrdb.service
[Unit]
Description=xrdb

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb %h/.Xresources
